# deploy/soulchain-interval.service
[Unit]
Description=SoulChain — Interval Daemon (Encrypted)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Update this path for your installation
ExecStart=/usr/local/bin/soulchain start --mode interval
WorkingDirectory=%h/.soulchain
Environment=SOULCHAIN_CONFIG=%h/.soulchain/soulchain.config.json
Environment=SOULCHAIN_PRIVATE_KEY=${SOULCHAIN_PRIVATE_KEY}
Environment=SOULCHAIN_KEYSTORE=%h/.soulchain/keystore.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soulchain

[Install]
WantedBy=multi-user.target
