# deploy/soulchain-on-write.service
[Unit]
Description=SoulChain — On-Write Daemon (Encrypted)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Update these paths for your installation
ExecStart=/usr/local/bin/soulchain start --mode on-write
WorkingDirectory=%h/.soulchain
Environment=SOULCHAIN_CONFIG=%h/.soulchain/soulchain.config.json
# Provide your private key via an env file or environment variable:
Environment=SOULCHAIN_PRIVATE_KEY=${SOULCHAIN_PRIVATE_KEY}
# Optional: encrypted keystore for Ed25519 crypto
Environment=SOULCHAIN_KEYSTORE=%h/.soulchain/keystore.json
# Set SOULCHAIN_KEYSTORE_PASSWORD in your environment or systemd override
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soulchain

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
